# scripts/concatman.service.tmpl
[Unit]
Description=Concatman will conncat string

[Service]
User={{USER}}
WorkingDirectory={{WORKING_DIR}}
ExecStart={{WORKING_DIR}}/bin/concatman
Restart=always

[Install]
WantedBy=multi-user.target